# docker/supervisord-azure.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/appcontrol/supervisord.log
pidfile=/var/run/supervisord.pid
childlogdir=/var/log/appcontrol

[program:gateway]
command=/usr/local/bin/appcontrol-gateway
user=appcontrol
autostart=true
autorestart=true
startretries=10
startsecs=3
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=RUST_LOG="%(ENV_RUST_LOG)s",GATEWAY_ID="%(ENV_GATEWAY_ID)s",GATEWAY_ZONE="%(ENV_GATEWAY_ZONE)s",LISTEN_ADDR="0.0.0.0",LISTEN_PORT="4443",BACKEND_URL="%(ENV_BACKEND_URL)s"

[program:agent]
command=/usr/local/bin/agent-wrapper.sh
user=root
autostart=true
autorestart=true
startretries=10
startsecs=5
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=RUST_LOG="%(ENV_RUST_LOG)s",HOME="/home/appcontrol",AGENT_ENROLLMENT_TOKEN="%(ENV_AGENT_ENROLLMENT_TOKEN)s",GATEWAY_ZONE="%(ENV_GATEWAY_ZONE)s"
